# copy/tmp/kumquat.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='application/kumquat' type='service' version='0'>
    <dependency name='net-physical' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/physical'/>
    </dependency>
    <dependency name='filesystem' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <exec_method name='start' type='method' exec='/opt/kumquat/manage.py backend' timeout_seconds='0'/>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='10'/>
    <exec_method name='refresh' type='method' exec=':kill -HUP' timeout_seconds='60'/>
    <property_group name='startd' type='framework'>
      <propval name='duration' type='astring' value='child'/>
      <propval name='ignore_error' type='astring' value='core,signal'/>
    </property_group>
    <instance name='default' enabled='false'>
      <method_context working_directory='/tmp'>
        <method_credential user='root' group='kumquat' />
        <method_environment>
          <envvar name='PATH' value='/opt/local/bin:/opt/local/sbin:/usr/bin:/usr/sbin:/bin'/>
          <envvar name='PYTHONPATH' value='/opt/local/lib/python2.7/site-packages'/>
        </method_environment>
      </method_context>
    </instance>
    <instance name='backend' enabled='false'>
      <method_context working_directory='/opt/kumquat'>
        <method_credential user='root' group='kumquat' />
        <method_environment>
          <envvar name='PATH' value='/opt/local/bin:/opt/local/sbin:/usr/bin:/usr/sbin:/bin'/>
        </method_environment>
      </method_context>
    </instance>
    <stability value='Evolving'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>Kumquat Backend Service</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
